# bmt_ai_os/runtime/init.d/bmt-firstboot
#!/sbin/openrc-run
# BMT AI OS — First-boot Initialization (OpenRC)
# Runs once to create data directories, pull default models, and set up configs.
# Safe to re-run: resumes from where it left off if interrupted.
#
# BMTOS-17 | Epic: BMTOS-EPIC-3 (OS Foundation & Infrastructure)

description="BMT AI OS first-boot initialization"

FIRSTBOOT_FLAG="/data/.firstboot-done"
DEFAULT_MODEL="${DEFAULT_MODEL:-qwen2.5-coder:7b-instruct-q4_K_M}"
OLLAMA_HOST="${OLLAMA_HOST:-http://127.0.0.1:11434}"

depend() {
    need docker
    need net
    after containerd docker bmt-recovery
    before ai-stack
}

start() {
    if [ -f "${FIRSTBOOT_FLAG}" ]; then
        einfo "First-boot already completed"
        return 0
    fi

    ebegin "Running BMT AI OS first-boot initialization"

    _create_data_dirs || { eend 1; return 1; }
    _init_configs     || { eend 1; return 1; }
    _pull_default_model
    _finalize

    eend 0
}

# ---------------------------------------------------------------------------
# Step 1: Create data directories on /data partition
# ---------------------------------------------------------------------------
_create_data_dirs() {
    einfo "Creating data directories on /data partition"

    local dirs="
        /data/ollama/models
        /data/chromadb
        /data/docker
        /data/bmt-ai-os
        /data/logs/bmt-ai-os
        /data/training/datasets
        /data/training/checkpoints
        /data/training/exports
    "

    local dir
    for dir in ${dirs}; do
        if [ ! -d "${dir}" ]; then
            mkdir -p "${dir}" || {
                eerror "Failed to create ${dir}"
                return 1
            }
        fi
    done

    # Set ownership for service users
    chmod 0755 /data/ollama /data/chromadb /data/docker
    einfo "Data directories created"
}

# ---------------------------------------------------------------------------
# Step 2: Initialize default configuration files
# ---------------------------------------------------------------------------
_init_configs() {
    einfo "Initializing default configuration"

    # Create controller config if absent
    local config_dir="/etc/bmt-ai-os"
    mkdir -p "${config_dir}"

    if [ ! -f "${config_dir}/controller.conf" ]; then
        cat > "${config_dir}/controller.conf" <<-CONF
# BMT AI OS Controller Configuration
# Generated by first-boot on $(date -u +%Y-%m-%dT%H:%M:%SZ)

OLLAMA_HOST=http://127.0.0.1:11434
CHROMADB_HOST=http://127.0.0.1:8000
DEFAULT_MODEL=${DEFAULT_MODEL}
LOG_LEVEL=info
CONF
    fi

    # Copy compose file if not already present
    if [ ! -f "${config_dir}/ai-stack/docker-compose.yml" ]; then
        mkdir -p "${config_dir}/ai-stack"
        if [ -f /usr/share/bmt-ai-os/ai-stack/docker-compose.yml ]; then
            cp /usr/share/bmt-ai-os/ai-stack/docker-compose.yml \
               "${config_dir}/ai-stack/docker-compose.yml"
        fi
    fi

    einfo "Configuration initialized"
}

# ---------------------------------------------------------------------------
# Step 3: Pull default model via Ollama
# ---------------------------------------------------------------------------
_pull_default_model() {
    # Flag to track model-pull progress so we can resume
    local model_flag="/data/.firstboot-model-pulled"

    if [ -f "${model_flag}" ]; then
        einfo "Default model already pulled"
        return 0
    fi

    einfo "Pulling default model: ${DEFAULT_MODEL}"
    einfo "This may take a while on the first boot..."

    # Start Ollama temporarily if it is not running
    local started_ollama=false
    if ! docker ps --format '{{.Names}}' 2>/dev/null | grep -q "bmt-ollama"; then
        docker run -d --name bmt-ollama-firstboot \
            -p 11434:11434 \
            -v ollama_models:/root/.ollama \
            ollama/ollama:latest \
            >>/var/log/bmt-ai-os/firstboot.log 2>&1
        started_ollama=true

        # Wait for Ollama API
        local wait=0
        while [ ${wait} -lt 60 ]; do
            if wget -q -O /dev/null "${OLLAMA_HOST}/api/tags" 2>/dev/null; then
                break
            fi
            sleep 2
            wait=$((wait + 2))
        done
    fi

    # Pull model (Ollama pull is resumable on its own)
    if wget -q -O /dev/null "${OLLAMA_HOST}/api/tags" 2>/dev/null; then
        wget --post-data "{\"name\":\"${DEFAULT_MODEL}\"}" \
            --header="Content-Type: application/json" \
            -q -O /dev/null \
            "${OLLAMA_HOST}/api/pull" 2>>/var/log/bmt-ai-os/firstboot.log
        if [ $? -eq 0 ]; then
            touch "${model_flag}"
            einfo "Default model pulled successfully"
        else
            ewarn "Model pull failed — will retry on next boot"
        fi
    else
        ewarn "Ollama not reachable — skipping model pull (will retry on next boot)"
    fi

    # Clean up temporary container
    if ${started_ollama}; then
        docker rm -f bmt-ollama-firstboot >>/var/log/bmt-ai-os/firstboot.log 2>&1
    fi
}

# ---------------------------------------------------------------------------
# Step 4: Finalize
# ---------------------------------------------------------------------------
_finalize() {
    touch "${FIRSTBOOT_FLAG}"
    einfo "First-boot initialization complete"
}
